# systemd/wazuh-weekly-report.timer
[Unit]
Description=Run Wazuh weekly report every Sunday 18:00

[Timer]
OnCalendar=Sun *-*-* 18:00:00
Persistent=true

[Install]
WantedBy=timers.target
